# ai-trading-server.service
[Unit]
Description=AI Trading Server
After=network.target

[Service]
Type=simple
User=mafuuux
WorkingDirectory=/home/mafuuux/ai-trading-server
Environment="PATH=/home/mafuuux/ai-trading-server/venv/bin"
Environment="ADMIN_USER=admin"
Environment="ADMIN_PASS=change_this_password"
Environment="SESSION_SECRET=change_this_secret"
ExecStart=/home/mafuuux/ai-trading-server/venv/bin/python /home/mafuuux/ai-trading-server/server.py
Restart=always
RestartSec=5
TimeoutStopSec=30
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
